# scripts/templates/systemd.service.tmpl
[Unit]
Description=@@INSTANCE@@ (secrets via Doppler)
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/docker/@@INSTANCE@@
Environment=HOME=/root
ExecStart=/usr/bin/doppler run -- /usr/bin/docker compose up -d
ExecStop=/usr/bin/docker compose -f /docker/@@INSTANCE@@/docker-compose.yml down

[Install]
WantedBy=multi-user.target
